# vendor/github.com/NVIDIA/gpu-monitoring-tools/exporters/prometheus-dcgm/bare-metal/prometheus-dcgm.service
[Unit]
Description=Prometheus DCGM exporter
Wants=dcgm.service
After=dcgm.service

[Service]
Type=simple
Restart=always
ExecStart=/usr/local/bin/dcgm-exporter

[Install]
WantedBy=multi-user.target
